# systemd/codex-telegram-bridge.service
[Unit]
Description=Codex Telegram Bridge
After=network-online.target
Wants=network-online.target
ConditionPathExists=%h/dev/telegram-agent-bridges/bridges/codex/.env

[Service]
Type=simple
WorkingDirectory=%h/dev/telegram-agent-bridges/bridges/codex
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/dev/telegram-agent-bridges/bridges/codex/run-bridge.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
